# com.sentinel-logger.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sentinel-logger.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alonmozes/qa tool/venv/bin/python</string>
        <string>/Users/alonmozes/qa tool/main.py</string>
        <string>--no-browser</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>9898</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/alonmozes/qa tool</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/alonmozes/qa tool/logs/sentinel-logger.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/alonmozes/qa tool/logs/sentinel-logger.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alonmozes/qa tool/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
